# modules/teleterm/files/teleterm.service
[Unit]
Description = runs teleterm
After=network.target

[Service]
ExecStart=/usr/bin/tt server
Restart=always

[Install]
WantedBy=default.target
