# modules/partofme/files/99-media-persistent.rules
SUBSYSTEM=="block" DM_NAME="partofme-data" RUN+="mount /media/persistent" RUN+="systemctl start pihole"
