# modules/partofme/files/99-media-persistent.rules
SUBSYSTEM=="block" ATTR{dm/name}="partofme-data" RUN+="/usr/bin/mount /media/persistent" RUN+="/usr/bin/systemctl start pihole"
